# container_files/tier-support/supervisord-tomee.conf
[program:tomee]
user=tomcat
command=/opt/tomee/bin/catalina.sh run 
stderr_logfile = /tmp/logtomcat
stderr_logfile_maxbytes=0
stdout_logfile = /tmp/logtomcat
stdout_logfile_maxbytes=0


