# container_files/tier-support/supervisord-loader.conf
[program:grouperloader]
user=shibd
command=bin/gsh.sh -loader
stderr_logfile = /tmp/loggrouper
stderr_logfile_maxbytes=0
stdout_logfile = /tmp/loggrouper
stdout_logfile_maxbytes=0
